# asset-tracker.service
[Unit]
Description=Personal Asset Tracker
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/asset-tracker
ExecStart=/path/to/asset-tracker/start_app.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
